# files2/etc/crontabs/root
*/1 * * * * pidof mwan3-watchdog || /usr/bin/mwan3-watchdog
